# autostart/Clipboard Manager.desktop
[Desktop Entry]
Exec=/PATH/TO/clipboardmanager
Name=Clipboard Manager
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
